# units/ersatz-key-unseal@.target
[Unit]
Description=External encrypted partition (%I) has unsealed everything
Documentation=https://github.com/moretiles/ersatz-key
