# units/unseal-usb.target
[Unit]
Description=USB and drive have unsealed everything
